# fastrpc-1.0.2/sensorspd
#!/sbin/openrc-run

description="sensorspd aDSP RPC daemon"

command="/usr/bin/adsprpcd"
command_args="sensorspd"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need localmount
    before iio-sensor-proxy
}

start_pre() {
    # Kondisi setara ConditionPathExists=|
    if [ ! -e /dev/fastrpc-adsp ] && [ ! -e /dev/fastrpc-adsp-secure ]; then
        eerror "fastrpc device not found"
        return 1
    fi
}
